# crous-bot.service
# ============================================================
#  Service systemd — pour l'option serveur 24h/24
#  (Oracle Cloud Always Free, VPS, Raspberry Pi...)
#
#  Installation (voir README, section "Option serveur") :
#    sudo cp deploy/crous-bot.service /etc/systemd/system/
#    sudo systemctl daemon-reload
#    sudo systemctl enable --now crous-bot
#
#  Suivre les logs en direct :
#    journalctl -u crous-bot -f
#
#  ⚠️ Adapte les deux chemins ci-dessous si ton dossier n'est pas
#     /home/ubuntu/crous-limoges-bot (utilisateur "ubuntu").
# ============================================================
[Unit]
Description=Bot de surveillance des logements CROUS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/crous-limoges-bot
ExecStart=/home/ubuntu/crous-limoges-bot/.venv/bin/python crous_bot.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
